# addons/spamhaus/cron.d/spamhausupdate
# ## SPAMHAUS (Adjust the value of the path with that of the customPath parameter)
30 19 * * * /usr/sbin/generatespamhaus.sh
