# addons/torbulkexit/cron.d/torlistupdate
# ## TOR (Adjust the value of the path with that of the customPath parameter)
30 19 * * * root /usr/sbin/generatetorbulkexit.sh >/dev/null 2>&1
